# thedevelopingdev/gitlab/packer-docker-setup/backups.cron
0,30 * * * * docker exec -t root_web_1 gitlab-backup create BACKUP=dump
0,30 * * * * cp -a /srv/gitlab/config /backup/gitlab-config
